# docker/supervisor/supervisord.conf
[supervisord]
logfile=supervisor.log
logfile_maxbytes=0
loglevel=info
pidfile=/var/supervisor/supervisord.pid
childlogdir=/var/supervisor/
stdout_logfile=/var/supervisor/supervisor.log
stdout_logfile_maxbytes=0

[unix_http_server]
file=/var/supervisor/supervisor.sock
chmod=0700

[inet_http_server]
port=9999
username=%(ENV_SUPERVISOR_API_USERNAME)s
password=%(ENV_SUPERVISOR_API_PASSWORD)s

[supervisorctl]
serverurl=unix:///var/supervisor/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[program:bankruptcy-worker]
process_name=%(process_num)02d
command=php /var/www/zaymer/src/laravel/artisan queue:work --queue=bankruptcy --sleep=20 --timeout=360 --tries=1
autostart=true
autorestart=true
user=%(ENV_SUPERVISOR_USER)s
numprocs=%(ENV_NUMPROCS_BANKRUPTCY_WORKER)s
redirect_stderr=true
stdout_logfile=/var/supervisor/supervisor.log
stdout_logfile_maxbytes=0

[program:underwriting-worker]
process_name=%(process_num)02d
command=php /var/www/zaymer/src/laravel/artisan queue:work --queue=underwriting --sleep=2 --timeout=60 --tries=2
autostart=true
autorestart=true
user=%(ENV_SUPERVISOR_USER)s
numprocs=%(ENV_NUMPROCS_UNDERWRITING_WORKER)s
redirect_stderr=true
stdout_logfile=/var/supervisor/supervisor.log
stdout_logfile_maxbytes=0

[program:potential-fraud-worker]
process_name=%(process_num)02d
command=php /var/www/zaymer/src/laravel/artisan queue:work --queue=potential_frauds --sleep=2 --timeout=360 --tries=1
autostart=true
autorestart=true
user=%(ENV_SUPERVISOR_USER)s
numprocs=%(ENV_NUMPROCS_POTENTIAL_FRAUD_WORKER)s
redirect_stderr=true
stdout_logfile=/var/supervisor/supervisor.log
stdout_logfile_maxbytes=0

[program:fail-remove-infinity-worker]
process_name=%(process_num)02d
command=php /var/www/zaymer/src/laravel/artisan queue:work --queue=infinity_remove --sleep=2 --timeout=10 --tries=2
autostart=true
autorestart=true
user=%(ENV_SUPERVISOR_USER)s
numprocs=%(ENV_NUMPROCS_FAIL_REMOVE_INFINITY_WORKER)s
redirect_stderr=true
stdout_logfile=/var/supervisor/supervisor.log
stdout_logfile_maxbytes=0

[program:triggers-worker]
process_name=%(process_num)02d
command=php /var/www/zaymer/src/laravel/artisan queue:work --queue=triggers --sleep=2 --timeout=360 --tries=1
autostart=true
autorestart=true
user=%(ENV_SUPERVISOR_USER)s
numprocs=%(ENV_NUMPROCS_TRIGGERS_WORKER)s
redirect_stderr=true
stdout_logfile=/var/supervisor/supervisor.log
stdout_logfile_maxbytes=0

[program:to-laravel-worker]
process_name=%(process_num)02d
command=php /var/www/zaymer/src/laravel/artisan worker:forwardable_job --count=500
autostart=true
autorestart=true
user=%(ENV_SUPERVISOR_USER)s
numprocs=%(ENV_NUMPROCS_TO_LARAVEL_WORKER)s
redirect_stderr=true
stdout_logfile=/var/supervisor/supervisor.log
stdout_logfile_maxbytes=0

[program:notification-worker]
process_name=%(process_num)02d
command=php /var/www/zaymer/src/laravel/artisan queue:work --queue=notification --sleep=2 --timeout=85 --tries=2
autostart=true
autorestart=true
user=%(ENV_SUPERVISOR_USER)s
numprocs=%(ENV_NUMPROCS_NOTIFICATION_WORKER)s
redirect_stderr=true
stdout_logfile=/var/supervisor/supervisor.log
stdout_logfile_maxbytes=0

[program:mass-notification-worker]
process_name=%(process_num)02d
command=php /var/www/zaymer/src/laravel/artisan queue:work --queue=mass_notification --sleep=2 --timeout=60 --tries=2
autostart=true
autorestart=true
user=%(ENV_SUPERVISOR_USER)s
numprocs=%(ENV_NUMPROCS_MASS_NOTIFICATION_WORKER)s
redirect_stderr=true
stdout_logfile=/var/supervisor/supervisor.log
stdout_logfile_maxbytes=0

[program:debt-notification-worker]
process_name=%(process_num)02d
command=php /var/www/zaymer/src/laravel/artisan queue:work --queue=debt_notification --sleep=2 --timeout=60 --tries=2
autostart=true
autorestart=true
user=%(ENV_SUPERVISOR_USER)s
numprocs=%(ENV_NUMPROCS_DEBT_NOTIFICATION_WORKER)s
redirect_stderr=true
stdout_logfile=/var/log/debt-notification-worker.log
stdout_logfile_maxbytes=0

[program:onec-owa-worker]
process_name=%(process_num)02d
command=php /var/www/zaymer/src/laravel/artisan one_c:workers:operator_work_application --count=500
autostart=true
autorestart=true
user=%(ENV_SUPERVISOR_USER)s
numprocs=%(ENV_NUMPROCS_ONEC_OWA_WORKER)s
redirect_stderr=true
stdout_logfile=/var/supervisor/supervisor.log
stdout_logfile_maxbytes=0

[program:discount-worker]
process_name=%(process_num)02d
command=php /var/www/zaymer/src/laravel/artisan queue:work --queue=discount --sleep=2 --timeout=60 --tries=2
autostart=true
autorestart=true
user=%(ENV_SUPERVISOR_USER)s
numprocs=%(ENV_NUMPROCS_DISCOUNT_WORKER)s
redirect_stderr=true
stdout_logfile=/var/supervisor/supervisor.log
stdout_logfile_maxbytes=0

[program:sms-notification-worker]
process_name=%(process_num)02d
command=php /var/www/zaymer/src/laravel/artisan queue:work --queue=sms_notification --sleep=2 --timeout=360 --tries=1
autostart=true
autorestart=true
user=%(ENV_SUPERVISOR_USER)s
numprocs=%(ENV_NUMPROCS_SMS_NOTIFICATION_WORKER)s
redirect_stderr=true
stdout_logfile=/var/supervisor/supervisor.log
stdout_logfile_maxbytes=0

[program:sms-priority-notification-worker]
process_name=%(process_num)02d
command=php /var/www/zaymer/src/laravel/artisan queue:work --queue=sms_priority_notification --sleep=2 --timeout=360 --tries=1
autostart=true
autorestart=true
user=%(ENV_SUPERVISOR_USER)s
numprocs=%(ENV_NUMPROCS_SMS_PRIORITY_NOTIFICATION_WORKER)s
redirect_stderr=true
stdout_logfile=/var/supervisor/supervisor.log
stdout_logfile_maxbytes=0

[program:change-customer-state-worker]
process_name=%(process_num)02d
command=php /var/www/zaymer/src/laravel/artisan queue:work --queue=change_customer_state --sleep=2 --timeout=360 --tries=1
autostart=true
autorestart=true
user=%(ENV_SUPERVISOR_USER)s
numprocs=%(ENV_NUMPROCS_CHANGE_CUSTOMER_STATE_WORKER)s
redirect_stderr=true
stdout_logfile=/var/supervisor/supervisor.log
stdout_logfile_maxbytes=0

[program:approve_identitication-worker]
process_name=%(process_num)02d
command=php /var/www/zaymer/src/laravel/artisan queue:work --queue=approve_identification --sleep=2 --timeout=360 --tries=1
autostart=true
autorestart=true
user=%(ENV_SUPERVISOR_USER)s
numprocs=%(ENV_NUMPROCS_APPROVE_IDENTIFICATION_WORKER)s
redirect_stderr=true
stdout_logfile=/var/supervisor/supervisor.log
stdout_logfile_maxbytes=0

[program:scoring-update-limit-worker]
process_name=%(process_num)02d
command=php /var/www/zaymer/src/laravel/artisan scoring:kafka_consumer --event-name=customerLimitUpdated
autostart=true
autorestart=true
user=%(ENV_SUPERVISOR_USER)s
numprocs=%(ENV_NUMPROCS_SCORING_UPDATE_LIMIT_WORKER)s
redirect_stderr=true
stdout_logfile=/var/supervisor/supervisor.log
stdout_logfile_maxbytes=0

[program:scoring-invalid-scoring-data-worker]
process_name=%(process_num)02d
command=php /var/www/zaymer/src/laravel/artisan scoring:kafka_consumer --event-name=invalidScoringData
autostart=true
autorestart=true
user=%(ENV_SUPERVISOR_USER)s
numprocs=%(ENV_NUMPROCS_SCORING_INVALID_SCORING_DATA_WORKER)s
redirect_stderr=true
stdout_logfile=/var/supervisor/supervisor.log
stdout_logfile_maxbytes=0

[program:scoring-customer-blocking-changed-worker]
process_name=%(process_num)02d
command=php /var/www/zaymer/src/laravel/artisan scoring:kafka_consumer --event-name=customerBlockingChanged
autostart=true
autorestart=true
user=%(ENV_SUPERVISOR_USER)s
numprocs=%(ENV_NUMPROCS_SCORING_CUSTOMER_BLOCKING_CHANGED_WORKER)s
redirect_stderr=true
stdout_logfile=/var/supervisor/supervisor.log
stdout_logfile_maxbytes=0

[program:scoring-customer-blacklist-changed-worker]
process_name=%(process_num)02d
command=php /var/www/zaymer/src/laravel/artisan scoring:kafka_consumer --event-name=customerBlacklistChanged
autostart=true
autorestart=true
user=%(ENV_SUPERVISOR_USER)s
numprocs=%(ENV_NUMPROCS_SCORING_CUSTOMER_BLACKLIST_CHANGED_WORKER)s
redirect_stderr=true
stdout_logfile=/var/supervisor/supervisor.log
stdout_logfile_maxbytes=0

[program:notifier-worker]
process_name=%(process_num)02d
command=php /var/www/zaymer/src/laravel/artisan queue:work --queue=notifier_queue --sleep=2 --timeout=360 --tries=1
autostart=true
autorestart=true
user=%(ENV_SUPERVISOR_USER)s
numprocs=%(ENV_NUMPROCS_NOTIFIER_WORKER)s
redirect_stderr=true
stdout_logfile=/var/supervisor/supervisor.log
stdout_logfile_maxbytes=0

[program:scan-passport-worker]
process_name=%(process_num)02d
command=php /var/www/zaymer/src/laravel/artisan workers:scan_passport --count=500
autostart=true
autorestart=true
user=%(ENV_SUPERVISOR_USER)s
numprocs=%(ENV_NUMPROCS_SCAN_PASSPORT_WORKER)s
redirect_stderr=true
stdout_logfile=/var/supervisor/supervisor.log
stdout_logfile_maxbytes=0

[program:onec-owc-worker]
process_name=%(process_num)02d
command=php /var/www/zaymer/src/laravel/artisan one_c:workers:operator_work_customer --count=500
autostart=true
autorestart=true
user=%(ENV_SUPERVISOR_USER)s
numprocs=%(ENV_NUMPROCS_ONEC_OWC_WORKER)s
redirect_stderr=true
stdout_logfile=/var/log/onec-owa-worker.log

[program:cron]
command=cron -f
redirect_stderr=true
stdout_logfile=/var/supervisor/supervisor.log
stdout_logfile_maxbytes=0

[program:send-operation-worker]
process_name=%(process_num)02d
command=php /var/www/zaymer/src/laravel/artisan queue:work --queue=send_operation --sleep=2 --timeout=60 --tries=2
autostart=true
autorestart=true
user=%(ENV_SUPERVISOR_USER)s
numprocs=%(ENV_NUMPROCS_SEND_OPERATION_WORKER)s
redirect_stderr=true
stdout_logfile=/var/supervisor/supervisor.log
stdout_logfile_maxbytes=0

[program:create-application-collector-discount-upload-worker]
process_name=%(process_num)02d
command=php /var/www/zaymer/src/laravel/artisan queue:work --queue=application_collector_discount_upload --sleep=2 --timeout=14400 --tries=1
autostart=true
autorestart=true
user=%(ENV_SUPERVISOR_USER)s
numprocs=%(ENV_NUMPROCS_CREATE_APPLICATION_COLLECTOR_DISCOUNT_UPLOAD_WORKER)s
redirect_stderr=true
stdout_logfile=/var/supervisor/supervisor.log
stdout_logfile_maxbytes=0

[program:receive_inn-worker]
process_name=%(process_num)02d
command=php /var/www/zaymer/src/laravel/artisan queue:work --queue=inn --sleep=2 --timeout=90 --tries=1
autostart=true
autorestart=true
user=%(ENV_SUPERVISOR_USER)s
numprocs=%(ENV_NUMPROCS_RECEIVE_INN_WORKER)s
redirect_stderr=true
stdout_logfile=/var/supervisor/supervisor.log
stdout_logfile_maxbytes=0

[program:decline-pixel-orders-worker]
process_name=%(process_num)02d
command=php /var/www/zaymer/src/laravel/artisan queue:work --queue=decline_pixel_orders --timeout=360 --tries=2
autostart=true
autorestart=true
user=%(ENV_SUPERVISOR_USER)s
numprocs=%(ENV_NUMPROCS_DECLINE_PIXEL_ORDERS_WORKER)s
redirect_stderr=true
stdout_logfile=/var/supervisor/supervisor.log
stdout_logfile_maxbytes=0

[program:rns_notification]
process_name=%(process_num)02d
command=php /var/www/zaymer/src/laravel/artisan queue:work --queue=rns_notification --sleep=2 --timeout=360
autostart=true
autorestart=true
user=%(ENV_SUPERVISOR_USER)s
numprocs=%(ENV_NUMPROCS_RNS_NOTIFICATION)s
redirect_stderr=true
stdout_logfile=/var/supervisor/supervisor.log
stdout_logfile_maxbytes=0

[program:create_notifications]
process_name=%(process_num)02d
command=php /var/www/zaymer/src/laravel/artisan queue:work --queue=create_notifications --sleep=2 --timeout=360
autostart=true
autorestart=true
user=%(ENV_SUPERVISOR_USER)s
numprocs=%(ENV_NUMPROCS_CREATE_NOTIFICATIONS)s
redirect_stderr=true
stdout_logfile=/var/supervisor/supervisor.log
stdout_logfile_maxbytes=0

[program:assign-collector-loans-upload]
process_name=%(process_num)02d
command=php /var/www/zaymer/src/laravel/artisan assign_collector_loans_upload:consume --count=500
autostart=true
autorestart=true
user=%(ENV_SUPERVISOR_USER)s
numprocs=%(ENV_NUMPROCS_ASSIGN_COLLECTOR_LOANS_UPLOAD)s
redirect_stderr=true
stdout_logfile=/var/supervisor/supervisor.log
stdout_logfile_maxbytes=0

[program:upload-application-to-collector-agency-worker]
process_name=%(process_num)02d
command=php /var/www/zaymer/src/laravel/artisan one_c:workers:collector_agency:handle_application --count=500
autostart=true
autorestart=true
user=%(ENV_SUPERVISOR_USER)s
numprocs=%(ENV_NUMPROCS_UPLOAD_APPLICATION_TO_COLLECTOR_AGENCY_WORKER)s
redirect_stderr=true
stdout_logfile=/var/supervisor/supervisor.log
stdout_logfile_maxbytes=0

[program:customer-duplicate]
process_name=%(process_num)02d
command=php /var/www/zaymer/src/laravel/artisan queue:work --queue=customer_duplicate --sleep=2 --timeout=60
autostart=true
autorestart=true
user=%(ENV_SUPERVISOR_USER)s
numprocs=%(ENV_NUMPROCS_CUSTOMER_DUPLICATE)s
redirect_stderr=true
stdout_logfile=/var/supervisor/supervisor.log
stdout_logfile_maxbytes=0

[program:billing-loan-updated-consume]
process_name=%(process_num)02d
command=php /var/www/zaymer/src/laravel/artisan billing:kafka_consumer_event --event-name=loanUpdated --group=zaymer-site-consumer-group
autostart=true
autorestart=true
user=%(ENV_SUPERVISOR_USER)s
numprocs=%(ENV_NUMPROCS_BILLING_LOAN_UPDATED_CONSUME)s
redirect_stderr=true
stdout_logfile=/var/supervisor/supervisor.log
stdout_logfile_maxbytes=0

[program:billing-msg-processing-error-consume]
process_name=%(process_num)02d
command=php /var/www/zaymer/src/laravel/artisan billing:kafka_consumer_event --event-name=msgProcessingError --group=zaymer-site-consumer-group
autostart=true
autorestart=true
user=%(ENV_SUPERVISOR_USER)s
numprocs=%(ENV_NUMPROCS_BILLING_MSG_PROCESSING_ERROR_CONSUME)s
redirect_stderr=true
stdout_logfile=/var/supervisor/supervisor.log
stdout_logfile_maxbytes=0

[program:infinity-sales-repeat-clients-third-sort]
process_name=%(process_num)02d
command=php /var/www/zaymer/src/laravel/artisan queue:work --queue=infinity_sales_repeat_clients_third_sort --sleep=2 --timeout=60
autostart=true
autorestart=true
user=%(ENV_SUPERVISOR_USER)s
numprocs=%(ENV_NUMPROCS_INFINITY_SALES_REPEAT_CLIENTS_THIRD_SORT)s
redirect_stderr=true
stdout_logfile=/var/supervisor/supervisor.log
stdout_logfile_maxbytes=0

[program:send-customer-action]
process_name=%(process_num)02d
command=php /var/www/zaymer/src/laravel/artisan queue:work --queue=send_customer_action --sleep=2 --timeout=60
autostart=true
autorestart=true
user=%(ENV_SUPERVISOR_USER)s
numprocs=%(ENV_NUMPROCS_SEND_CUSTOMER_ACTION)s
redirect_stderr=true
stdout_logfile=/var/supervisor/supervisor.log
stdout_logfile_maxbytes=0

[program:billing-worker]
process_name=%(process_num)02d
command=php /var/www/zaymer/src/laravel/artisan queue:work --queue=billing --sleep=2 --timeout=60
autostart=true
autorestart=true
user=%(ENV_SUPERVISOR_USER)s
numprocs=%(ENV_NUMPROCS_BILLING_WORKER)s
redirect_stderr=true
stdout_logfile=/var/supervisor/supervisor.log
stdout_logfile_maxbytes=0

[program:change-card-scoring-result-worker]
process_name=%(process_num)02d
command=php /var/www/zaymer/src/laravel/artisan queue:work --queue=change_card_scoring_result --sleep=2 --timeout=360 --tries=1
autostart=true
autorestart=true
user=%(ENV_SUPERVISOR_USER)s
numprocs=%(ENV_NUMPROCS_CARD_SCORING_WORKER)s
redirect_stderr=true
stdout_logfile=/var/supervisor/supervisor.log
stdout_logfile_maxbytes=0

[program:request-inn-worker]
process_name=%(process_num)02d
command=php /var/www/zaymer/src/laravel/artisan queue:work --queue=request_inn --sleep=2 --timeout=60
autostart=true
autorestart=true
user=%(ENV_SUPERVISOR_USER)s
numprocs=%(ENV_NUMPROCS_REQUEST_INN_WORKER)s
redirect_stderr=true
stdout_logfile=/var/supervisor/supervisor.log
stdout_logfile_maxbytes=0
